# meta-agl-profile-core/recipes-multimedia/lightmediascanner/files/lightmediascanner.service
[Unit]
Description=Lightmediascanner (LMS)

[Service]
Type=dbus
BusName=org.lightmediascanner
UMask=0077
SmackProcessLabel=User::App::agl-service-mediascanner
ExecStart=/usr/bin/lightmediascannerd --startup-scan --db-path=/run/user/%U/usrshr/lightmediascannerd.db --directory=/media --directory=/home/0/Music --directory=/home/0/Videos

[Install]
WantedBy=multi-user.target
Alias=dbus-org.lightmediascanner.service
