# meta-agl/recipes-multimedia/lightmediascanner/files/90-lightmediascanner.rules
SUBSYSTEM=="block", ACTION=="add"    RUN+="/etc/udev/scripts/lightmediascanner.sh"
